# usbnotify/00-usb.rules
KERNEL=="sd*",RUN+="/opt/usbnotify.sh"

